# common/root/etc/s6-overlay/s6-rc.d/svc-dcv/run
#!/usr/bin/with-contenv bash

set -e

# Only run in DCV mode
if [ "${REMOTE_PROTOCOL}" != "dcv" ]; then
	echo "svc-dcv: REMOTE_PROTOCOL=${REMOTE_PROTOCOL}, exiting (DCV disabled)"
	exec sleep infinity
fi

# Wait for the X server on ${DISPLAY} to be ready before starting DCV.
# s6 "started" for svc-xvfb only means the run script was exec'd, not that
# the X display exists yet. If dcvserver starts first, its console-session
# agent crash-loops on the missing display and clients that connect during
# that window get "No session is available or you are not authorised".
echo "Waiting for X display ${DISPLAY} to be ready..."
while ! xset q &>/dev/null; do
	sleep 0.5
done
echo "X display ${DISPLAY} is ready."

echo "Starting DCV Server..."

mkdir -p /run/dcv /run/dcvlogon

DCV_ARGS="-d --web-port ${HTTP_PORT:-3000} --auth none --create-session --session-type=console --display ${DISPLAY}${DCV_LICENSE_FILE:+ --license-file ${DCV_LICENSE_FILE}}"

exec /usr/bin/dcvserver ${DCV_ARGS}
