# common/root/etc/s6-overlay/s6-rc.d/svc-dcv/run
#!/usr/bin/with-contenv bash

set -e

# Only run in DCV mode
if [ "${REMOTE_PROTOCOL}" != "dcv" ]; then
	echo "svc-dcv: REMOTE_PROTOCOL=${REMOTE_PROTOCOL}, exiting (DCV disabled)"
	exec sleep infinity
fi

# Wait for the X server on ${DISPLAY} to be ready before starting DCV.
# s6 "started" for svc-xvfb only means the run script was exec'd, not that
# the X display exists yet. If dcvserver starts first, its console-session
# agent crash-loops on the missing display and clients that connect during
# that window get "No session is available or you are not authorised".
echo "Waiting for X display ${DISPLAY} to be ready..."
while ! xset q &>/dev/null; do
	sleep 0.5
done
echo "X display ${DISPLAY} is ready."

# Surface DCV logs to container output. dcvserver only logs to files, so
# session/agent failures (the reason a client sees "No session is available")
# are invisible in `docker logs` without this.
DCV_LOG_DIR="${HOME}/.local/share/NICE/dcv/log"
mkdir -p "${DCV_LOG_DIR}"
pkill -f "tail .*NICE/dcv/log" 2>/dev/null || true
tail -n +1 -F \
	"${DCV_LOG_DIR}/server.log" \
	"${DCV_LOG_DIR}/agent.root.console.log" 2>/dev/null &

# Handle a restart: the /usr/bin/dcvserver wrapper does not exec the real
# binary, so a previous instance survives service restart and keeps the web
# port bound. Clean up before starting.
pkill -f "dcv/dcvserver" 2>/dev/null && sleep 2 || true
pkill -f "dcv/dcvagent" 2>/dev/null || true

# dcvserver launched outside systemd reads ${HOME}/.config/dcv/dcv.conf and
# ignores /etc/dcv/dcv.conf — mirror the config generated by init-dcv-config
# to where this process will actually read it.
mkdir -p "${HOME}/.config/dcv"
cp -f /etc/dcv/dcv.conf "${HOME}/.config/dcv/dcv.conf"

# Audio. Desktop apps play into a null sink and DCV grabs its monitor —
# the same sinks Selkies mode creates (svc-selkies is disabled in DCV mode,
# so create them here). Wait for PulseAudio the same way svc-selkies does.
until [ -f "${PULSE_RUNTIME_PATH}/pid" ]; do
	sleep .5
done
if [ ! -f /dev/shm/audio.lock ]; then
	s6-setuidgid "${USER}" pactl load-module module-null-sink \
		sink_name="output" sink_properties=device.description="output"
	s6-setuidgid "${USER}" pactl load-module module-null-sink \
		sink_name="input" sink_properties=device.description="input"
	touch /dev/shm/audio.lock
fi
s6-setuidgid "${USER}" pactl set-default-sink output

# Microphone. The agent creates a virtual mic by loading module-pipe-source
# with a FIFO under ${HOME}/.cache/NICE/dcv. The agent runs as root and
# would create that directory root-owned — the user's PulseAudio daemon
# then cannot create the FIFO inside it and the mic reports Unavailable.
# Pre-create it as the user.
s6-setuidgid "${USER}" mkdir -p "${HOME}/.cache/NICE/dcv"

# Recording apps should default to the DCV virtual microphone once the
# agent creates it — otherwise the first source (a sink monitor) is the
# default and apps record the desktop instead of the client mic.
(
	until s6-setuidgid "${USER}" pactl list short sources 2>/dev/null | grep -q "AWS-Virtual-Microphone"; do
		sleep 2
	done
	s6-setuidgid "${USER}" pactl set-default-source "AWS-Virtual-Microphone"
) &

# The DCV server/agent run as root while PulseAudio is the user's. A root
# libpulse client that discovers the server via the shared
# PULSE_RUNTIME_PATH "secures" that directory by chowning it to root,
# which locks every user application out of audio (apps then fail
# libpulse's ownership check and play nothing — the desktop goes silent
# while DCV faithfully streams an idle sink). Point the DCV processes at
# the socket directly and give them a private runtime dir instead.
export PULSE_SERVER="unix:${PULSE_RUNTIME_PATH}/native"
export PULSE_RUNTIME_PATH="/tmp/pulse-dcv"

echo "Starting DCV Server..."

mkdir -p /run/dcv /run/dcvlogon

# Internal port only — nginx serves ${PREFIX} on 3000 and proxies to it
# (same architecture as Selkies with its 8081 backend, see svc-nginx)
DCV_ARGS="-d --web-port 8082 --auth none --create-session --session-type=console --display ${DISPLAY}${DCV_LICENSE_FILE:+ --license-file ${DCV_LICENSE_FILE}}"

exec /usr/bin/dcvserver ${DCV_ARGS}
