# common/root/etc/s6-overlay/s6-rc.d/svc-dcv/run
#!/usr/bin/with-contenv bash

set -e

# Only run in DCV mode
if [ "${REMOTE_PROTOCOL}" != "dcv" ]; then
	echo "svc-dcv: REMOTE_PROTOCOL=${REMOTE_PROTOCOL}, exiting (DCV disabled)"
	exec sleep infinity
fi

DCV_SESSION_TYPE="${DCV_SESSION_TYPE:-console}"

# Surface DCV logs to container output. dcvserver only logs to files, so
# session/agent failures (the reason a client sees "No session is available")
# are invisible in `docker logs` without this.
DCV_LOG_DIR="${HOME}/.local/share/NICE/dcv/log"
mkdir -p "${DCV_LOG_DIR}"
pkill -f "tail .*NICE/dcv/log" 2>/dev/null || true
tail -n +1 -F \
	"${DCV_LOG_DIR}/server.log" \
	"${DCV_LOG_DIR}/agent.root.console.log" \
	"${DCV_LOG_DIR}/agent.${USER}.helios.log" 2>/dev/null &

# Handle a restart: the /usr/bin/dcvserver wrapper does not exec the real
# binary, so a previous instance survives service restart and keeps the web
# port bound. Clean up before starting.
pkill -f "dcv/dcvserver" 2>/dev/null && sleep 2 || true
pkill -f "dcv/dcvagent" 2>/dev/null || true
pkill -f "dcv/dcvsessionlauncher" 2>/dev/null || true

# dcvserver launched outside systemd reads ${HOME}/.config/dcv/dcv.conf and
# ignores /etc/dcv/dcv.conf — mirror the config generated by init-dcv-config
# to where this process will actually read it.
mkdir -p "${HOME}/.config/dcv"
cp -f /etc/dcv/dcv.conf "${HOME}/.config/dcv/dcv.conf"

mkdir -p /run/dcv /run/dcvlogon

# Internal port only — nginx serves ${PREFIX} on 3000 and proxies to it
# (same architecture as Selkies with its 8081 backend, see svc-nginx)
DCV_ARGS="-d --web-port 8082 --auth none${DCV_LICENSE_FILE:+ --license-file ${DCV_LICENSE_FILE}}"

if [ "${DCV_SESSION_TYPE}" = "virtual" ]; then
	# Virtual session: DCV owns the X server (Xdcv) instead of scraping the
	# Xvfb display — no Xvfb/svc-desktop, dynamic display sizing, cheaper
	# frame capture. How this works without extra container capabilities:
	#
	# - dcvserver runs as ${USER}, not root. The server authorizes its agent
	#   by reading /proc/<pid>/exe, which across different uids would need
	#   CAP_SYS_PTRACE; with server and agent both running as ${USER} the
	#   kernel allows it with no added capability.
	# - dcvsessionlauncher spawns the session (Xdcv + init script) and must
	#   stay root: it uses runuser to become ${USER}.
	# - The two find each other on a private system-type D-Bus instance
	#   (/run/dcv/bus): the launcher hard-codes the system bus while a
	#   non-systemd dcvserver autolaunches a session bus, so neither default
	#   bus lets them meet. The launcher exits if the server's bus name is
	#   absent, so it is started only after the server registers.
	echo "Starting DCV in virtual session mode..."

	export DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dcv/bus"
	export DBUS_SYSTEM_BUS_ADDRESS="unix:path=/run/dcv/bus"
	if [ ! -S /run/dcv/bus ]; then
		dbus-daemon --config-file=/opt/helios/dcv-bus.conf --fork --nopidfile
	fi

	# The session and the server both run as ${USER} and need to own their
	# state; the root-run console mode may have left root-owned files behind.
	chown -R "${USER}:${USER}" /run/dcv /run/dcvlogon
	mkdir -p "${HOME}/.local" "${HOME}/.config" "${HOME}/.cache"
	chown -R "${USER}:${USER}" "${HOME}/.local" "${HOME}/.config" "${HOME}/.cache"

	# Once the server registers on the bus: start the launcher, then create
	# the session; recreate it if it ever drops.
	(
		until dbus-send --system --print-reply --dest=org.freedesktop.DBus \
			/ org.freedesktop.DBus.ListNames 2>/dev/null |
			grep -q com.nicesoftware.DcvServer; do
			sleep 1
		done
		/usr/bin/dcvsessionlauncher &
		sleep 2
		while true; do
			if ! s6-setuidgid "${USER}" dcv list-sessions 2>/dev/null | grep -q "Session:"; then
				s6-setuidgid "${USER}" dcv create-session \
					--type=virtual \
					--owner "${USER}" \
					--init /opt/helios/dcv-startwm.sh \
					helios 2>/dev/null || true
			fi
			sleep 5
		done
	) &

	echo "Starting DCV Server..."
	exec s6-setuidgid "${USER}" /usr/bin/dcvserver ${DCV_ARGS}
fi

# Console session: attach to the Xvfb display. Wait for the X server on
# ${DISPLAY} first — s6 "started" for svc-xvfb only means the run script
# was exec'd, not that the display exists yet. If dcvserver starts first,
# its console-session agent crash-loops on the missing display and clients
# that connect during that window get "No session is available".
echo "Waiting for X display ${DISPLAY} to be ready..."
while ! xset q &>/dev/null; do
	sleep 0.5
done
echo "X display ${DISPLAY} is ready."

echo "Starting DCV Server..."

exec /usr/bin/dcvserver ${DCV_ARGS} --create-session --session-type=console --display "${DISPLAY}"
